# pkg/xarph-desktop/usr/lib/systemd/user/xarph-session.target
[Unit]
Description=Xarph Desktop Session
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target
RefuseManualStart=yes
RefuseManualStop=yes

[Install]
WantedBy=graphical-session.target
